# files/etc/supervisor.d/pdns-recursor.ini
[program:pdns-recursor]
command=/usr/sbin/pdns_recursor
	--allow-from=0.0.0.0/0,::/0
	--any-to-tcp=yes
	--client-tcp-timeout=2
	--daemon=no
	--disable-packetcache=no
	--disable-syslog=yes
	--dnssec=process-no-validate
	--dnssec-log-bogus=no
	--dont-query=0.0.0.0/8,10.0.0.0/8,100.64.0.0/10,127.0.0.0/8,169.254.0.0/16,172.16.0.0/12,192.0.0.0/24,192.0.2.0/24,192.168.0.0/16,198.18.0.0/15,198.51.100.0/24,203.0.113.0/24,224.0.0.0/3,::1/128,fc00::/7,fe80::/10,::/96,::ffff:0:0/96,100::/64,2001:db8::/32
	--edns-outgoing-bufsize=1680
	--edns-subnet-whitelist=255.255.255.255/32
	--entropy-source=/dev/urandom
	--etc-hosts-file=/dev/null
	--export-etc-hosts=no
	--hint-file=/etc/pdns/opennic.hint
	--latency-statistic-size=10000
	--local-address=0.0.0.0,::
	--local-port=53
	--non-local-bind=yes
	--loglevel=6
	--log-common-errors=yes
	--lowercase-outgoing=no
	--max-cache-entries=32768
	--max-cache-ttl=604800
	--max-mthreads=2048
	--max-packetcache-entries=32768
	--max-qperq=50
	--max-negative-ttl=900
	--max-recursion-depth=40
	--max-tcp-clients=1024
	--max-tcp-per-client=64
	--max-total-msec=15000
	--minimum-ttl-override=30
	--network-timeout=1000
	--packetcache-ttl=3600
	--packetcache-servfail-ttl=60
	--pdns-distributes-queries=no
	--query-local-address=0.0.0.0
	--query-local-address6=::
	--quiet=yes
	--reuseport=yes
	--root-nx-trust=yes
	--security-poll-suffix=
	--serve-rfc1918=yes
	--server-down-max-fails=16
	--server-down-throttle-time=60
	--setgid=recursor
	--setuid=recursor
	--single-socket=no
	--spoof-nearmiss-max=10
	--stack-size=200000
	--stats-ringbuffer-entries=10000
	--threads=8
	--trace=no
	--udp-truncation-threshold=512
	--use-incoming-edns-subnet=no
redirect_stderr=true
