# files/etc/supervisor.d/pdns-recursor.ini
[program:pdns-recursor]
command=/usr/sbin/pdns_recursor
	--allow-from=%(ENV_ALLOW_FROM)s
	--any-to-tcp=%(ENV_ANY_TO_TCP)s
	--client-tcp-timeout=%(ENV_CLIENT_TCP_TIMEOUT)s
	--daemon=%(ENV_DAEMON)s
	--disable-packetcache=%(ENV_DISABLE_PACKETCACHE)s
	--disable-syslog=%(ENV_DISABLE_SYSLOG)s
	--dnssec=%(ENV_DNSSEC)s
	--dnssec-log-bogus=%(ENV_DNSSEC_LOG_BOGUS)s
	--dont-query=%(ENV_DONT_QUERY)s
	--edns-outgoing-bufsize=%(ENV_EDNS_OUTGOING_BUFSIZE)s
	--edns-subnet-whitelist=%(ENV_EDNS_SUBNET_WHITELIST)s
	--entropy-source=%(ENV_ENTROPY_SOURCE)s
	--etc-hosts-file=%(ENV_ETC_HOSTS_FILE)s
	--export-etc-hosts=%(ENV_EXPORT_ETC_HOSTS)s
	--forward-zones-file=%(ENV_FORWARD_ZONES_FILE)s
	--hint-file=%(ENV_HINT_FILE)s
	--latency-statistic-size=%(ENV_LATENCY_STATISTIC_SIZE)s
	--local-address=%(ENV_LOCAL_ADDRESS)s,%(ENV_HEALTHCHECK_LISTEN_IP)s
	--local-port=%(ENV_LOCAL_PORT)s
	--non-local-bind=%(ENV_NON_LOCAL_BIND)s
	--loglevel=%(ENV_LOGLEVEL)s
	--log-common-errors=%(ENV_LOG_COMMON_ERRORS)s
	--lowercase-outgoing=%(ENV_LOWERCASE_OUTGOING)s
	--max-cache-entries=%(ENV_MAX_CACHE_ENTRIES)s
	--max-cache-ttl=%(ENV_MAX_CACHE_TTL)s
	--max-mthreads=%(ENV_MAX_MTHREADS)s
	--max-packetcache-entries=%(ENV_MAX_PACKETCACHE_ENTRIES)s
	--max-qperq=%(ENV_MAX_QPERQ)s
	--max-negative-ttl=%(ENV_MAX_NEGATIVE_TTL)s
	--max-recursion-depth=%(ENV_MAX_RECURSION_DEPTH)s
	--max-tcp-clients=%(ENV_MAX_TCP_CLIENTS)s
	--max-tcp-per-client=%(ENV_MAX_TCP_PER_CLIENT)s
	--max-total-msec=%(ENV_MAX_TOTAL_MSEC)s
	--minimum-ttl-override=%(ENV_MINIMUM_TTL_OVERRIDE)s
	--network-timeout=%(ENV_NETWORK_TIMEOUT)s
	--packetcache-ttl=%(ENV_PACKETCACHE_TTL)s
	--packetcache-servfail-ttl=%(ENV_PACKETCACHE_SERVFAIL_TTL)s
	--pdns-distributes-queries=%(ENV_PDNS_DISTRIBUTES_QUERIES)s
	--query-local-address=%(ENV_QUERY_LOCAL_ADDRESS)s
	--query-local-address6=%(ENV_QUERY_LOCAL_ADDRESS6)s
	--quiet=%(ENV_QUIET)s
	--reuseport=%(ENV_REUSEPORT)s
	--root-nx-trust=%(ENV_ROOT_NX_TRUST)s
	--security-poll-suffix=%(ENV_SECURITY_POLL_SUFFIX)s
	--serve-rfc1918=%(ENV_SERVE_RFC1918)s
	--server-down-max-fails=%(ENV_SERVER_DOWN_MAX_FAILS)s
	--server-down-throttle-time=%(ENV_SERVER_DOWN_THROTTLE_TIME)s
	--setgid=%(ENV_SETGID)s
	--setuid=%(ENV_SETUID)s
	--single-socket=%(ENV_SINGLE_SOCKET)s
	--spoof-nearmiss-max=%(ENV_SPOOF_NEARMISS_MAX)s
	--stack-size=%(ENV_STACK_SIZE)s
	--stats-ringbuffer-entries=%(ENV_STATS_RINGBUFFER_ENTRIES)s
	--threads=%(ENV_THREADS)s
	--trace=%(ENV_TRACE)s
	--udp-truncation-threshold=%(ENV_UDP_TRUNCATION_THRESHOLD)s
	--use-incoming-edns-subnet=%(ENV_USE_INCOMING_EDNS_SUBNET)s
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true
